# .config/systemd/user/journal-clean.service
[Unit]
Description=Clean old journal entries

[Service]
Type=oneshot
ExecStart=/usr/bin/journalctl --vacuum-time=2weeks
