# .config/systemd/user/xresources.path
[Unit]
Description=Xresources changes monitoring

[Path]
PathModified=/home/den/.Xresources

[Install]
WantedBy=multi-user.target
