# .config/systemd/user/xresources.service
[Unit]
Description=Xresources changes monitoring

[Service]
Type=oneshot
ExecStart=systemctl --user restart dunst ; polybar-msg cmd restart
